# tinyprogrammer.service
[Unit]
Description=TinyProgrammer - Autonomous Code Writing Device
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/aerovisual/TinyProgrammer
# Hide console cursor before starting
ExecStartPre=/bin/sh -c 'echo 0 > /sys/class/graphics/fbcon/cursor_blink'
ExecStartPre=-/bin/sh -c 'setterm -cursor off > /dev/tty1 2>/dev/null || true'
ExecStart=/usr/bin/python3 /home/aerovisual/TinyProgrammer/main.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/tinyprogrammer.log
StandardError=append:/var/log/tinyprogrammer.log

[Install]
WantedBy=multi-user.target
